# supervisor_ble_server.conf
[program:ble_server]
command=bash -c 'sleep 5 && python /home/greenponik/ble_server/main.py'
autostart=true
autorestart=true
